# visitors_addon/files/php83/etc/services.d/php-fpm83/run
#!/usr/bin/with-contenv bashio

mkdir -p /tmp/sessions
chown nginx:nginx /tmp/sessions

exec php-fpm83 --nodaemonize
