# visitors_addon/files/php83/etc/services.d/php-fpm83/run
#!/bin/execlineb -S0
with-contenv
php-fpm83 -F
