# visitors_addon/files/general/etc/services.d/nginx/finish
#!/bin/execlineb -S0
foreground { /bin/rm -f /tmp/.container_finished }
/bin/s6-svscanctl -t /var/run/s6/services
