# visitors_addon/files/general/etc/services.d/nginx/run
#!/bin/execlineb -S0
with-contenv
nginx
