# visitors_addon/files/general/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::net.wait_for 9000

exec nginx
